# supervisord-cita-checker.conf
[program:cita-checker]
command=/bin/bash -lc 'until [ -x /home/nonroot/venv/bin/python3 ]; do sleep 2; done; until /home/nonroot/venv/bin/python3 -c "import seleniumbase" >/dev/null 2>&1; do sleep 2; done; export DISPLAY=:1; exec /home/nonroot/venv/bin/python3 /home/nonroot/cita-checker.py'
autostart=true
autorestart=true
startsecs=5
startretries=10
user=nonroot
stdout_logfile=/tmp/cita-checker.out
stderr_logfile=/tmp/cita-checker.err
priority=40
